# roles/alerta/templates/alerta-mailer.service.j2
[Unit]
Description=alerta-mailer

[Service]
Type=simple
ExecStart={{alerta_server_dir}}/bin/alerta-mailer

[Install]
WantedBy=multi-user.target
