# roles/alerta/templates/alertad.service.j2
[Unit]
Description=alertad

[Service]
Type=simple
ExecStart={{alerta_server_dir}}/bin/alertad

[Install]
WantedBy=multi-user.target
